# contrib/cirros/resources/etc/service/consul-template/run
#!/bin/bash

consul-template \
  -config=/etc/consul-template/consul-template.conf
